# spring-boot-crud-redis.service
[Unit]
Description=Spring Boot CRUD Redis Application
After=network.target

[Service]
Type=simple
User=deployer
WorkingDirectory=/home/deployer/spring-boot-crud-redis
ExecStart=/usr/bin/java -jar /home/deployer/spring-boot-crud-redis/target/crud-redis-0.0.1-SNAPSHOT.jar
SuccessExitStatus=143
Restart=always
RestartSec=10

# Load environment variables from .env file
EnvironmentFile=/home/deployer/spring-boot-crud-redis/.env

# Environment variables
Environment="JAVA_OPTS=-Xmx512m -Xms256m"
Environment="SPRING_PROFILES_ACTIVE=dev"

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=spring-boot-crud-redis

[Install]
WantedBy=multi-user.target
